# crontab
* * * * * curl https://www.groupmuse.com/cron/test?api_key=aadf04e0b039c26611719009824e79ab

# An empty line is required at the end of this file for a valid cron file.
