# crontab
* * * * * curl -k https://www.groupmuse.com/cron/test?api_key=$CRON_API_KEY

# An empty line is required at the end of this file for a valid cron file.
